# metadata/supervisord/metad/metad.ini
[program:metad]


command=start-metad.sh

stdout_logfile=/opt/openpitrix/log/metad.log
stdout_logfile_maxbytes=500MB
stdout_logfile_backups=10
redirect_stderr=true

user=root
autorestart=true
